# crontab/changtou.sh
#!/usr/bin/env bash
cd /root/brzhang/crontab

/usr/local/bin/node changtou.js
